# pi-gen/stage2.1/00-ad-hoc-net/files/batman-start.service
[Unit]
Description=Start batman interface
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/local/sbin/mesh.sh

[Install]
WantedBy=multi-user.target
